# connection.service
[Unit]
Description = dckr
After = docker.service

[Service]
ExecStart=/bin/bash /usr/bin/containers.sh

[Install]
WantedBy = multi-user.target
